# nfo2web.cron
# Update the hackbox-system-monitor interface hourly
@hourly root	hackbox-system-monitor
